# rootfs/etc/services.d/insync/finish
#!/usr/bin/with-contenv bash

insync quit
