# rootfs/etc/services.d/insync/finish
#!/command/with-contenv /bin/sh

insync quit
