# rootfs/etc/services.d/insync/run
#!/usr/bin/with-contenv bash

echo '[INSYNC] Starting Insync.'
insync start --no-daemon
