# rootfs/etc/services.d/insync/run
#!/command/with-contenv /bin/sh

echo '[INSYNC] Starting Insync.'
insync start --no-daemon
